# files/etc/systemd/system/zuul-scheduler.service
[Unit]
Description=Zuul Scheduler Service
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/bin/sh -c "${PREFIX}/bin/zuul-scheduler -d"
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
